# files/supervisord/proxy-8080.conf
[program:proxy8080]
command=/usr/bin/kubectl --kubeconfig="/root/.kube/config" port-forward service/ambassador 8080:80
user=root
