# files/supervisord/proxy-8500.conf
[program:proxy8500]
command=/usr/bin/kubectl --kubeconfig="/root/.kube/config" port-forward service/consul-server 8501:8500
user=root
startretries=10
