# files/supervisord/proxy-8877.conf
[program:proxy8877]
command=/usr/local/bin/proxy-8877.sh
user=root
