# files/supervisord/proxy-8080.conf
[program:proxy8080]
command=/usr/local/bin/proxy-8080.sh
user=root
startretries=10
